# vsphere/CENTOS7_OS_POST_BUILD/iscsi-auto-login.service
[Unit]
Description=iSCSI Auto-Login Service
After=network.target

[Service]
Type=oneshot
ExecStart=/sbin/iscsiadm --mode node --loginall=all
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
